# shared/cron/crontab
* * * * * cd /var/www/html/backend && php artisan schedule:run >> /dev/stdout 2>&1